# debian/evocheck.cron.d
33 1 2-31 * * root /usr/share/scripts/evocheck.sh --cron
33 1 1    * * root /usr/share/scripts/evocheck.sh
